# linux/etc/systemd/system/aggregator-collector-elma-vm-access-doc.timer
[Unit]
Description=Elma VmAccessDoc collector timer

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=aggregator.target