# linux/etc/systemd/system/aggregator-collector-storages.timer
[Unit]
Description=Storages collector timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=aggregator.target