# linux/etc/systemd/system/aggregator.target
[Unit]
Description=Aggregator timers target
Requires=aggregator-collector-hosts.timer
After=timers.target

[Install]
WantedBy=multi-user.target