# linux/etc/systemd/system/aggregator-collector-hosts.timer
[Unit]
Description=Hosts collector timer

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=aggregator.target